# roles/guifibages_server/files/etc/profile.d/guifibages.sh
if [ -f /guifibages/etc/profile ]; then
  . /guifibages/etc/profile
fi
